# Library/LaunchAgents/com.local.border-cycle.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.local.border-cycle</string>

    <!-- パスは chezmoi が機械別に解決（ユーザー名非依存）。 -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>{{ .chezmoi.homeDir }}/.config/borders/border-cycle</string>
    </array>

    <!-- borders(Homebrew) を解決。launchd 既定 PATH には無いため必須。 -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>{{ .chezmoi.homeDir }}/.local/state/border-cycle.out.log</string>

    <key>StandardErrorPath</key>
    <string>{{ .chezmoi.homeDir }}/.local/state/border-cycle.err.log</string>
</dict>
</plist>
